# chef-repo/site-cookbooks/granagile/templates/default/etc/cron.d/create-alminium-vcs.erb
MAILTO=""
PATH=/usr/local/bin:/usr/bin:/bin
* * * * * root sh /opt/granagile/script/create-alminium-git.sh >> /tmp/create-alminium-git.log 2>&1
* * * * * root sh /opt/granagile/script/create-alminium-svn.sh >> /tmp/create-alminium-svn.log 2>&1
