# chef-repo/site-cookbooks/centos/templates/default/etc/profile.d/git-completion.sh.erb
source /usr/share/git-core/contrib/completion/git-completion.bash
source /usr/share/git-core/contrib/completion/git-prompt.sh
export PS1='[\[\033[32m\]\u@\h\[\033[00m\] \[\033[34m\]\W\[\033[00m\]]\[\033[31m\]$(__git_ps1)\[\033[00m\]\$ '
